# dynamite_nsm/configurations/configurations/base/default_configs/systemd/dynamite-agent.target
[Unit]
Description=Dynamite Agent Service
Requires=multi-user.target
After=multi-user.target
Wants=suricata.service zeek.service filebeat.service


[Install]
WantedBy=multi-user.target