# dynamite_nsm/services/systemd/dynamite.service
[Unit]
Description=Dynamite Agent Service 
Requires=network.target pf_ring-ko.service
After=network.target pf_ring-ko.service 
