# dynamite_nsm/services/systemd/dynamite-scanner.target
[Unit]
Description=Dynamite Scanner Service 
Requires=network-online.target
After=network-online.target
Wants=rumble.service filebeat.service

